# config/sync/ai.openclaw.homelab-sync-hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.homelab-sync-hermes</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>git -C /Users/ryan/Developer/homelab fetch origin &amp;&amp; git -C /Users/ryan/Developer/homelab merge --ff-only origin/main || (echo "[homelab-sync-hermes] non-ff detected, halting" &gt;&amp;2; exit 1)</string>
  </array>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/homelab-sync-hermes.log</string>
</dict>
</plist>
